# custom-stage/00-sensos-hotspot/files/sensos-hotspot.service
[Unit]
Description=SensOS Wi-Fi hotspot
Wants=NetworkManager.service
After=NetworkManager.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/sensos-start-hotspot
ExecStop=/usr/local/sbin/sensos-stop-hotspot
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
